# pyload-ng/rootfs/etc/s6-overlay/s6-rc.d/discovery/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Home Assistant Community Add-on: pyLoad-ng
# Sends discovery information to Home Assistant.
# ==============================================================================
declare config
declare proto
# Wait for pyLoad to become available
bashio::net.wait_for 8000 127.0.0.1 300

if bashio::config.true ssl; then
    proto="https"
else
    proto="http"
fi


config=$(bashio::var.json \
    url "${proto}://$(bashio::addon.hostname):$(bashio::addon.port 8000)/" \
    username "pyload" \
    password "pyload" \
)

if bashio::discovery "pyload" "${config}" > /dev/null; then
    bashio::log.info "Successfully sent discovery information to Home Assistant."
else
    bashio::log.error "Discovery message to Home Assistant failed!"
fi
